# build/X86-nas/files/etc/crontabs/root
* * * * * sleep 10s; /mnt/nvme0n1p1/opt/stop_hassio_audio.sh > /dev/null 2>&1
